# apps/links/devfiles/docker/rootfs/build/services/node.conf
[program:node]

command = "/build/bin/start_node"
autostart = true
autorestart = true